# obos-genissue.service
[Unit]
Description=Generate dynamic /run/issue
DefaultDependencies=no
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/bin/gen-issue
