# obos-genissue.timer
[Unit]
Description=Refresh /run/issue every 5 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
Unit=gen-issue.service

[Install]
WantedBy=timers.target
